# packaging/init/systemd/unitylan-engine.service
[Unit]
Description=UnityLAN mesh engine
Documentation=https://github.com/gtosh4/unity_lan
# The mesh dials peers over the internet; wait for real connectivity.
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/lib/unitylan/unitylan-engine run /etc/unitylan/engine.toml
# `always` (not `on-failure`) so the signed auto-update path can exit(0) after swapping the binary
# and be relaunched onto the new version. An explicit `systemctl stop` still stops (not a restart).
Restart=always
RestartSec=5
# Teardown reverts the firewall and resolver first, then the interface — and that last step can wedge
# inside boringtun's uapi. Cap the wait well under systemd's 90s default so a stuck `remove_interface`
# doesn't hold a stop (or a reboot) for a minute and a half; the host-global state is already undone
# by then, and a leftover interface/uapi socket is recovered on the next start.
TimeoutStopSec=15

# Run as root (User unset) but with the `unitylan` group, so RuntimeDirectory (/run/unitylan) and
# the control socket the engine binds there are group-owned by `unitylan`. Members of that group
# (the desktop user, added at install) reach the socket without root — no chown needed, so no
# CAP_CHOWN in the bounding set below. The group is created by the package postinstall.
Group=unitylan

# /var/lib/unitylan (0700 root:unitylan — group has no access) — WG private key + pinned trust
# anchor. State survives restarts.
StateDirectory=unitylan
StateDirectoryMode=0700
# /run/unitylan — root:unitylan 0750 so `unitylan` group members can traverse to the control socket.
RuntimeDirectory=unitylan
RuntimeDirectoryMode=0750

# Privilege: userspace WireGuard (ip/TUN), the nftables firewall, and resolvectl all need
# CAP_NET_ADMIN; the `.internal` resolver may bind :53 (CAP_NET_BIND_SERVICE). Nothing else.
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_BIND_SERVICE
DeviceAllow=/dev/net/tun rw

# Hardening — kept compatible with netlink (WG), the nft pipe, and the resolved D-Bus call.
NoNewPrivileges=yes
ProtectSystem=strict
# The signed auto-update rewrites the engine binary in place, so its dir must stay writable despite
# ProtectSystem=strict (which otherwise makes /usr read-only). The binary lives in its own dir (not
# /usr/bin) so this grant is scoped to the engine alone — not every binary in /usr/bin (e.g. sudo).
ReadWritePaths=/usr/lib/unitylan
ProtectHome=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
MemoryDenyWriteExecute=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
